# config/rubber/common/remote_syslog.upstart.conf
<%
  @path = "/etc/init/remote_syslog.conf"
  @backup = false
%>
description "Monitor files and send to remote syslog"
start on runlevel [2345]
stop on runlevel [!2345]

respawn

pre-start exec /usr/bin/test -e /etc/log_files.yml

exec /usr/local/rubies/1.9.3-p392/bin/remote_syslog -D --tls
